# data/systemd/pyfpmd.service
[Unit]
Description=Pamac

[Service]
Type=dbus
BusName=org.frugalware.pyfpm
ExecStart=/usr/bin/pyfpmd
